# scripts/service/smbd/run
#!/bin/sh

set -eu

cd /run/samba/
exec 2>&1
exec /usr/sbin/smbd --foreground --debug-stdout --debuglevel=3
